# bootalert.service
[Unit]
Description=Bootalert
After=postfix.service

[Service]
ExecStart=/usr/local/bin/bootalert.py

[Install]
WantedBy=multi-user.target

